# scripts/dev.nzt108.github-radar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.nzt108.github-radar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nzt108/Projects/obsidian-second-mind/.venv/bin/python</string>
        <string>/Users/nzt108/Projects/obsidian-second-mind/scripts/github_radar_cron.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Run daily at 10:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/github-radar.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/github-radar.log</string>

    <!-- Don't run if missed (e.g. laptop was sleeping) — run on next wake -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
</dict>
</plist>
